# com.claudecode.screenshot-uploader.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claudecode.screenshot-uploader</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>export PATH=/opt/homebrew/bin:$PATH; ~/claude-screenshot-uploader/screenshot_uploader.sh</string>
    </array>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <true/>
    
    <key>StandardOutPath</key>
    <string>/tmp/screenshot-uploader.log</string>
    
    <key>StandardErrorPath</key>
    <string>/tmp/screenshot-uploader-error.log</string>
    
    <key>WorkingDirectory</key>
    <string>~/</string>
    
    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>